# solidblocks-ansible/blcks/utils/roles/md_monitor/templates/md-monitor.service.j2
[Unit]
Description=md monitor

[Service]
Type=oneshot
ExecStart={{ bin_dir }}/md-monitor.sh
Environment="MD_MONITOR_LOG_FILE={{ logfile }}"

[Install]
WantedBy=multi-user.target