# solidblocks-ansible/blcks/utils/roles/md_monitor/templates/md-monitor.timer.j2
[Unit]
Description=md monitor

[Timer]
OnCalendar={{ schedule }}
Unit=md-monitor.service

[Install]
WantedBy=multi-user.target
