# solidblocks-rds-postgresql-ansible/ansible/blcks/rds_postgresql/roles/pgbackrest/templates/backup-incr.timer.j2
[Unit]
Description={{ stanza_name }} database incr backup

[Timer]
OnCalendar={{backup_incr_schedule}}
Unit={{ stanza_name }}-backup-incr.service

[Install]
WantedBy=multi-user.target
