# solidblocks-rds-postgresql-ansible/ansible/blcks/rds_postgresql/roles/pgbackrest/templates/backup-full.timer.j2
[Unit]
Description={{ stanza_name }} database full backup

[Timer]
OnCalendar={{backup_full_schedule}}
Unit={{ stanza_name }}-backup-full.service

[Install]
WantedBy=multi-user.target
