# solidblocks-rds-postgresql-ansible/ansible/blcks/rds_postgresql/roles/pgbackrest/templates/backup-full.service.j2
[Unit]
Description={{ pgbackrest_stanza_name }} database full backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/{{ pgbackrest_stanza_name }}-backup-full.sh

[Install]
WantedBy=multi-user.target