# solidblocks-rds-postgresql-ansible/ansible/blcks/rds_postgresql/roles/pgbackrest/templates/backup-incr.service.j2
[Unit]
Description={{ stanza_name }} database incr backup

[Service]
Type=oneshot
ExecStart=/usr/local/bin/{{ stanza_name }}-backup-incr.sh

[Install]
WantedBy=multi-user.target