# solidblocks-rds-postgresql-ansible/ansible/blcks/rds_postgresql/roles/common/templates/pgbackrest/backup-full.timer.j2
[Unit]
Description={{ pgbackrest_stanza_name }} database full backup

[Timer]
OnCalendar=*-*-* 20:00:00
Unit={{ pgbackrest_stanza_name }}-backup-full.service

[Install]
WantedBy=multi-user.target
